# scripts/vwait.service
[Unit]
Description=vagrant waiter

[Service]
Type=oneshot
ExecStart=/usr/bin/vwait

[Install]
WantedBy=multi-user.target
